# templates/default/teamcity-build-agent.conf.erb
# teamcity-build-agent - TeamCity Build Agent
# THIS FILE IS MANAGED BY CHEF. ALL CHANGES WILL BE OVERWRITTEN!
#
description "TeamCity Build Agent"

start on runlevel [2345]
stop on runlevel [016]

# TeamCity can take a while to shutdown
kill timeout 30

respawn

script

  export JRE_HOME=/usr

  exec start-stop-daemon --start -c www-data --exec "<%= @agent_path %>/bin/agent.sh" -- run

end script